# etc/cron.hourly/slus
#!/bin/sh

slus="/usr/bin/slus"

[ -x $slus] || exit 0

$slus > /dev/null

exit 0
